# system_config/meb.service
[Unit]
Description=med (MEb) Nuxt application service
Documentation=https://example.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meb
Group=meb
WorkingDirectory=/home/staging/med
Environment=NODE_ENV=production
# Updated PATH to include the node bin directory (not the executable itself)
Environment=PATH=/home/staging/.nvm/versions/node/v22.21.1/bin:/usr/bin:/usr/local/bin:/bin:/usr/sbin:/sbin
EnvironmentFile=/etc/default/staging

# Allow ample time for `npm install` / `npm run build` during startup
TimeoutStartSec=1200
TimeoutStopSec=30

# Build steps (run as the service user)
# Using absolute path to npm
ExecStartPre=/home/staging/.nvm/versions/node/v22.21.1/bin/npm install
ExecStartPre=/home/staging/.nvm/versions/node/v22.21.1/bin/npm run build

# Start the server
# Using absolute path to npm
ExecStart=/home/staging/.nvm/versions/node/v22.21.1/bin/npm run start:server

# Keep the service running and restart on failure
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging
SyslogIdentifier=meb
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
